# scripts/watcher.service
[Unit]
Description=YT Strategy Agent watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=l-admin
WorkingDirectory=/srv/ai-hub/workspaces/yt-strategy-agent
EnvironmentFile=/srv/ai-hub/workspaces/yt-strategy-agent/.env
ExecStart=/srv/ai-hub/workspaces/yt-strategy-agent/.venv/bin/python /srv/ai-hub/workspaces/yt-strategy-agent/watcher.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
